# husky_bringup/upstart/job.conf
description	"bringup robot-job"

start on net-device-up IFACE=interface0
stop on net-device-down IFACE=interface0

console output
respawn

exec robot-job-start
pre-stop exec robot-job-stop
